# .codeploy/portable.fs/usr/lib/systemd/system/ops-maria.service
[Unit]
Wants=maria@ops\x2dmaria.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c "echo provision ops-maria DB schema"
